# config-app/ds5-haptics-loopback.service
[Unit]
Description=DS5Dongle auto-haptics audio loopback (active only while the dongle is present)
Documentation=https://github.com/loteran/DS5Dongle
After=pipewire.service wireplumber.service pipewire-pulse.service
Wants=pipewire.service
# This unit is started on demand by udev (ENV{SYSTEMD_USER_WANTS}) when the
# DualSense dongle is plugged in, and stopped by udev when it is unplugged.
# It must NOT be enabled at boot: running it without the dongle present would
# make the playback stream fall back to the speakers and create an audio
# feedback loop. The ExecStartPre guard below is a second line of defense.
# Bound automatic restarts so a permanently broken state can't spin forever.
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
# Safety guard + boot/plug race handling: wait (up to ~15s) for the dongle sink
# (ds5_dongle_sink, renamed by 51-ds5dongle.conf) to actually appear before
# launching the loopback. At boot or plug-in, udev may start us before
# WirePlumber has applied the rename; testing only once and giving up left
# pw-loopback to exit immediately and haptics dead until the next manual plug.
# If the sink never shows up (dongle truly absent) the guard fails and the
# service does not start — no feedback loop.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do pactl list short sinks | grep -q "\\bds5_dongle_sink\\b" && exit 0; sleep 0.5; done; exit 1'
# Capture the monitor of the current default sink and pipe it to the dongle.
# node.dont-reconnect=true on the playback side ensures that if the dongle
# disappears at runtime, the stream stays disconnected instead of falling
# back to the default output (no feedback loop).
ExecStart=/usr/bin/pw-loopback -c 2 \
  --capture-props="node.name=ds5_haptics_capture media.class=Stream/Input/Audio stream.capture.sink=true target.object=@DEFAULT_AUDIO_SINK@ audio.position=[FL,FR] stream.dont-remix=true node.passive=true" \
  --playback=ds5_dongle_sink \
  --playback-props="node.name=ds5_haptics_playback node.dont-reconnect=true audio.position=[AUX0,AUX1] stream.dont-remix=true node.passive=true"
# Restart on any unexpected exit (pw-loopback should run as long as the dongle
# is present). systemd does NOT restart after a commanded stop, so the udev
# unplug handler still stops it cleanly with no feedback risk.
Restart=always
RestartSec=2
